# Tugas4/case5/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
pidfile=/tmp/supervisord.pid
minfds=1024
minprocs=200
directory=/tmp
childlogdir=/tmp

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true
priority=100

[program:backend_migrate]
command=/app/hasil migrate
autorestart=false
priority=150

[program:backend_seed]
command=/app/hasil seed
autorestart=false
priority=160

[program:backend]
command=/app/hasil
autorestart=true
priority=200
